# systemd/ubuntu-aktualizacje.service
[Unit]
Description=Ubuntu_Aktualizacje — Full System Update
Documentation=https://github.com/KasprowiczM/Ubuntu_Aktualizacje
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Run as the user who owns the project (set via install-systemd-timer.sh)
User=%i
Group=%i
WorkingDirectory=/home/%i/Dev_Env/Ubuntu_Aktualizacje

# Environment
Environment=HOME=/home/%i
Environment=DISPLAY=:0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
Environment=XDG_RUNTIME_DIR=/run/user/1000

# Skip drivers on automated runs (require manual review for firmware)
ExecStart=/home/%i/Dev_Env/Ubuntu_Aktualizacje/update-all.sh --no-drivers

# Allow full runtime
TimeoutStartSec=3600

# Logging
StandardOutput=append:/home/%i/Dev_Env/Ubuntu_Aktualizacje/logs/systemd_update.log
StandardError=append:/home/%i/Dev_Env/Ubuntu_Aktualizacje/logs/systemd_update.log

[Install]
WantedBy=multi-user.target
