# systemd/ubuntu-aktualizacje@.service
[Unit]
Description=Ubuntu_Aktualizacje — Full System Update for %i
Documentation=https://github.com/KasprowiczM/Ubuntu_Aktualizacje
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=%i
Group=%i
WorkingDirectory=/home/%i/Dev_Env/Ubuntu_Aktualizacje

Environment=HOME=/home/%i
Environment=DISPLAY=:0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
Environment=XDG_RUNTIME_DIR=/run/user/1000

# Skips the explicit drivers/firmware module. APT may still update kernels,
# microcode, Mesa, and other standard OS packages according to apt policy.
ExecStart=/usr/bin/env bash /home/%i/Dev_Env/Ubuntu_Aktualizacje/update-all.sh --no-drivers

TimeoutStartSec=3600
StandardOutput=append:/home/%i/Dev_Env/Ubuntu_Aktualizacje/logs/systemd_update.log
StandardError=append:/home/%i/Dev_Env/Ubuntu_Aktualizacje/logs/systemd_update.log

[Install]
WantedBy=multi-user.target
